# trunk/embox/third-party/emlinux/tools/inittab
# System initialization (runs when system boots).

::sysinit:/etc/init.d/rc.S
::respawn:/bin/sh
