# tags/emlinux-0.1/usr/etc/inittab
# System initialization (runs when system boots).

::sysinit:/etc/rc.S
::respawn:/bin/sh
